# rtic/infra/systemd/rtic-ddns.timer
[Unit]
Description=rtic DDNSŽ•ľ 5Ž∂ĄŽßąŽč§ ÍįĪžč†

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
